# Meta/waechterd.service
[Unit]
Description=Wächter daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as root so privileged init can happen; daemon drops privileges internally.
User=root
Group=root

ExecStart=/usr/bin/waechterd
Restart=on-failure
RestartSec=2s

# If your daemon uses SIGTERM cleanly, this is fine.
KillSignal=SIGTERM
TimeoutStopSec=15s

# Hardening (keep conservative since it manipulates networking/bpf)
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
